# deploy/systemd/dumpster.service
[Unit]
Description=Dumpster Forum Server
Documentation=https://github.com/yourorg/dumpster
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=dumpster
Group=dumpster
WorkingDirectory=/opt/dumpster

# Environment configuration
EnvironmentFile=/opt/dumpster/.env

# Main process
ExecStart=/opt/dumpster/bin/dumpster
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dumpster

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Allow write access to specific directories
ReadWritePaths=/opt/dumpster/tmp
ReadWritePaths=/opt/dumpster/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
