# deploy/systemd/dumpster-xf-chat.service
[Unit]
Description=Dumpster XenForo-Compatible Chat Server
Documentation=https://github.com/yourorg/dumpster
After=network.target
# Note: xf-chat connects to XenForo MySQL, not local PostgreSQL

[Service]
Type=simple
User=dumpster
Group=dumpster
WorkingDirectory=/opt/dumpster

# Environment configuration
EnvironmentFile=/opt/dumpster/.env
EnvironmentFile=/opt/dumpster/.env.xf-chat

# Main process
ExecStart=/opt/dumpster/bin/xf-chat
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dumpster-xf-chat

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Allow write access to specific directories
ReadWritePaths=/opt/dumpster/tmp
ReadWritePaths=/opt/dumpster/logs

# Resource limits (higher for chat server)
LimitNOFILE=131072
LimitNPROC=8192

[Install]
WantedBy=multi-user.target
